# root/etc/services.d/autosub/run
#!/usr/bin/with-contenv bash

cd /app/autosub || exit

exec \
	s6-setuidgid abc python AutoSub.py -c /config/config.properties -l
